# deploy/systemd/cloudpanel-gateway-helper.service
[Unit]
Description=CloudPanel Gateway privileged CLI helper
After=network.target cloudpanel-gateway-nginx-commit.service
Requires=cloudpanel-gateway-nginx-commit.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/cloudpanel-gateway --config /etc/cloudpanel-gateway/config.json helper
Restart=on-failure
RestartSec=3
RuntimeDirectory=cloudpanel-gateway
# Other users may traverse this directory, but only root and the gateway group
# can open the 0660 Unix socket created inside it.
RuntimeDirectoryMode=0755
# CloudPanel's own clpctl invokes sudo for several supported operations.
# The helper is root-only and accepts a validated Unix-socket protocol instead.
PrivateTmp=true
# clpctl creates and manages CloudPanel site directories under /home.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
